# root/etc/services.d/transcoder/run
#!/usr/bin/with-contenv bash

echo "Starting Transcoder"
exec \
  s6-setuidgid abc python /config/app/transcoder.py
